# systemd/engram-reconcile.service
[Unit]
Description=Engram reconciliation pipeline
# Allow up to 6 starts per hour. Soft-exits from the VRAM check count
# as successful exits (exit 0), but tight retry loops on fundamental
# failures should still be caught.
StartLimitBurst=6
StartLimitIntervalSec=3600

[Service]
Type=oneshot
ExecStart=%h/.cargo/bin/engram reconcile
# Reconcile can take 10+ minutes on large backlogs
TimeoutStartSec=1800
# Log level controlled by ~/.config/engram/config.toml (default: warn)
# Override here for debugging: Environment=RUST_LOG=engram=info
