# systemd/engram-debounce.timer
[Unit]
Description=Debounce timer for engram reconciliation

[Timer]
# Wait 15 minutes after the last file change before triggering reconcile.
# Each new .path activation resets the timer. Coarse interval is fine —
# engram's value is cross-session, not intra-session. Incremental
# processing means each run only handles new content.
OnActiveSec=15min
AccuracySec=1min
Unit=engram-reconcile.service
